# mouse/logid.timer
[Unit]
Description=Run logid on boot with delay
After=network.target

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
