# mouse/logid.timer
[Unit]
Description=Run logid on boot with delay

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target
